# etc/cron.d/spacewalk-sst-sync
# Sync data from spacewalk to candlepin every 4 hours
* */4 * * * root /usr/bin/spacewalk-sst-sync
