# scripts/spacewalk-splice-tool.cron
0 */4 * * * root /bin/bash --login -c "/usr/bin/spacewalk-splice-checkin > /dev/null"