# scripts/spacewalk-splice-tool.cron
0 */4 * * * root /usr/bin/spacewalk-splice-checkin
